# systemd/clearscan-bot.service
[Unit]
Description=ClearScan Telegram Bot
After=network.target

[Service]
Type=simple
User=clearscan
Group=clearscan
WorkingDirectory=/opt/clearscan
Environment=PYTHONPATH=/opt/clearscan
ExecStart={{PYTHON3_PATH}} /opt/clearscan/src/telegram_bot.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target